# erdi-crm-upload/ops/erdi-mail-pull.timer
[Unit]
Description=Run ERDI CRM Gmail import every five minutes

[Timer]
OnBootSec=2min
OnUnitInactiveSec=5min
Persistent=true
Unit=erdi-mail-pull.service

[Install]
WantedBy=timers.target
